# systemd/dbus/dbus.initd
#!/sbin/openrc-run

description="Fake dbus dep"

depend()
{
	:
}

start()
{
	:
}
